# xivbot.service
[Unit]
Description=xivbot server

[Service]
ExecStart=/etc/cqhttp/xivbot/xivbot
KillMode=process
Restart=on-failure